# conf/supervisord.conf
[supervisord]
nodaemon=true
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:usbmuxd]
command=usbmuxd -f -v
autostart=true
autorestart=true
stderr_logfile=/app/logs/usbmuxd.err.log
stdout_logfile=/app/logs/usbmuxd.out.log
priority=10


[program:netmuxd]
command=/app/bin/netmuxd
autostart=true
autorestart=true
stderr_logfile=/app/logs/netmuxd.err.log
stdout_logfile=/app/logs/netmuxd.out.log
priority=30
